# deploy/supervisord/uwsgi.conf
[program:uwsgi]
command=uwsgi ../uwsgi.ini
autostart=true
autorestart=true
stdout_logfile=/var/log/uwsgi/ablog.log
stderr_logfile=/var/log/uwsgi/ablog_err.log
redirect_stderr=true
stopsignal=QUIT
